# buildroot-external/rootfs-overlay/usr/lib/systemd/system/rauc.service.d/haos-rockpi.conf
# change tmp dir to allow download of large bundles
[Unit]
RequiresMountsFor=/mnt/data

[Service]
Environment=TMPDIR=/mnt/data/tmp
ExecStartPre=/usr/bin/mkdir -p /mnt/data/tmp
